# scripts/sandbox/systemd/team-alpha-nats.service
[Unit]
Description=team-alpha NATS bus (VM-local)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/nats-server -c /etc/team-alpha/nats.conf
Restart=on-failure
RestartSec=2

DynamicUser=yes
SupplementaryGroups=team-alpha
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
